# .config/autostart/autostart.desktop
[Desktop Entry]
Type=Application
Exec=kitty fish -c "neofetch; exec fish"
